# dist/rpm/manticore.service.in
[Unit]
Description=Manticore Search Engine
After=network.target remote-fs.target nss-lookup.target
After=syslog.target

[Service]
Type=notify
User=manticore
Group=manticore
# Run ExecStartPre with root-permissions
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p @CMAKE_INSTALL_FULL_RUNSTATEDIR@/manticore
ExecStartPre=/bin/chown manticore.manticore @CMAKE_INSTALL_FULL_RUNSTATEDIR@/manticore
# Run ExecStart with User=manticore / Group=manticore
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/searchd --config @CMAKE_INSTALL_FULL_SYSCONFDIR@/manticoresearch/manticore.conf --nodetach $_ADDITIONAL_SEARCHD_PARAMS
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
LimitCORE=infinity
LimitMEMLOCK=infinity
Restart=on-failure
TimeoutStopSec=60
NotifyAccess=all

[Install]
WantedBy=multi-user.target
Alias=manticore.service
Alias=searchd.service
